# etc/ontap.ellson.com.crontab
# $Id$ $Revision$

# SOURCE_DIR=CURRENT for daily snapshot builds
# or SOURCE_DIR=ARCHIVE for release builds
SOURCE_DIR=CURRENT
#SOURCE_DIR=ARCHIVE

15 0 * * * (cd graphviz-build; ../bin/anoncvs.tcl -qz3 update -d; crontab etc/ontap.ellson.com.crontab)
 0 2 * * * graphviz-build/redhat/graphviz-bin-rpm.tcl $SOURCE_DIR
 0 3 * * * graphviz-build/redhat/graphviz-bin-att-rpm.tcl $SOURCE_DIR
