# etc/penguin.ellson.com.crontab
# SOURCE_DIR=CURRENT for daily snapshot builds
# or SOURCE_DIR=ARCHIVE for release builds
SOURCE_DIR=CURRENT
#SOURCE_DIR=ARCHIVE

# Directory containing build scripts from graphviz-build CVS repository
BIN=graphviz-build/redhat

 0 2 * * * $BIN/graphviz-bin-rpm.tcl $SOURCE_DIR
 0 3 * * * $BIN/graphviz-bin-att-rpm.tcl $SOURCE_DIR
